# package/rooter/ext-rooter-basic/files/etc/init.d/ttl
#!/bin/sh /etc/rc.common

. /lib/functions.sh

START=98

start()
{
	/usr/lib/custom/ttl.sh &
}
